# submission/week-01/07.service
[Unit]
Description=DevOps Hello World Oneshot Service

[Service]
Type=oneshot
ExecStart=/bin/echo "Hello World"